# debian/ngcp-acc-cdi.cron.d
* * * * *   root    /usr/sbin/ngcp-acc-cdi >/dev/null
